# data/salt/states/base/salt/config/etc/cron.d/salt-provisionnew
*/10 * * * * root /root/scripts/solo.py -port=65001 /usr/bin/salt-ssh  '*' -i state.sls salt > /dev/null 2>&1
