# data/salt/states/base/salt/config/etc/cron.d/salt-batchstate
*/5 * * * * root /root/scripts/solo.py -port=65000 /usr/bin/salt -b 1 '*' state.highstate > /dev/null 2>&1
