# data/salt/states/runbook/redis/config/supervisord.conf
[supervisord]
nodaemon = true

[program:stunnel-server]
command = /usr/bin/stunnel4 /config/stunnel-server.conf
priority = 1
autostart = true
autorestart = true
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:redis]
command = /usr/bin/redis-server /config/redis.conf
priority = 2
autostart = true
autorestart = true
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
