# data/salt/states/runbook/runbook/monitors/worker/config/supervisord.conf
[supervisord]
nodaemon = true

[program:stunnel-client]
command = /usr/bin/stunnel4 /config/stunnel-client.conf
priority = 1
autostart = true
autorestart = true
startretries = 9999999
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:worker]
process_name = %(program_name)s_%(process_num)02d
command = /code/worker.py /config/worker.yml
directory = /code/
numproc = 5
priority = 2
autostart = true
autorestart = true
startretries = 9999999
user = runapp
group = users
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
