# data/salt/states/runbook/runbook/bridge/config/supervisord.conf
[supervisord]
nodaemon = true

[program:stunnel-client]
command = /usr/bin/stunnel4 /config/stunnel-client.conf
priority = 1
autostart = true
autorestart = true
startretries = 9999999
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:bridge]
command = /code/bridge.py /config/bridge.yml
directory = /code/
priority = 2
autostart = true
autorestart = true
startsecs = 60
startretries = 9999999
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
