# data/salt/states/runbook/supervisor/config/etc/cron.d/docker-kill
* * * * * root /root/scripts/docker-kill.sh 2>&1 > /dev/null
