# data/salt/states/runbook/runbook/web/config/supervisord.conf
[supervisord]
nodaemon = true

[program:stunnel-client]
command = /usr/bin/stunnel4 /config/stunnel-client.conf
priority = 1
autostart = true
autorestart = true
startretries = 9999999
user = root
group = root
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:nginx]
command = /usr/sbin/nginx
priority = 2
autostart = true
autorestart = true
startretries = 9999999
user = root
group = root
directory = /code
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:uwsgi]
command = /usr/bin/uwsgi --ini /config/uwsgi.cfg
autostart = true
autorestart = true
startretries = 9999999
user = runapp
group = users
directory = /code
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stopsignal=QUIT

[program:genstatic]
command = /usr/bin/python /code/genstatic.py
autostart = true
autorestart = unexpected
startretries = 999999
user = runapp
group = users
directory = /data/runstatic/public_html
umask = 022
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
